# server/docker-webhook.service
[Unit]
Description=Docker Update Webhook Server
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/scripts
ExecStart=/usr/bin/python3 /opt/scripts/webhook-server.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Variables de entorno si necesario
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
